# reference/prototype/cc-token-rotator.service
[Unit]
Description=Claude Max two-account OAuth token rotator (no-ops until accounts/ENABLED exists)

[Service]
Type=oneshot
ExecStart=%h/git/theconnman/claude-settings/services/token-rotator/rotate.sh
